# webui/harness-webui.service
[Unit]
Description=Multi-AI Harness WebUI
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ws/harness/webui
ExecStart=/home/ubuntu/ws/harness/webui/start.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
